# doc/core/howto/listings/systemd/twistedweb.socket
[Socket]
ListenStream=0.0.0.0:8080

[Install]
WantedBy=sockets.target
