# docs/core/howto/listings/systemd/www.example.com.socket
[Socket]
ListenStream=0.0.0.0:80
FileDescriptorName=my-web-port

[Install]
WantedBy=sockets.target
