# trunk/docs/core/howto/listings/systemd/www.example.com.socket
[Socket]
ListenStream=0.0.0.0:80

[Install]
WantedBy=sockets.target
